# deploy/playbooks/roles/app/files/app.service
[Unit]
Description=starter service
After=network.target

[Service]
Type=simple
User=starter
Group=starter
WorkingDirectory=/usr/local/starter
ExecStart=/usr/local/starter/app -conf /usr/local/starter/app.yml
ExecReload=/bin/kill -SIGUSR2 $MAINPID
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target